# tomcat7_jdk7_nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:tomcat]
command=catalina.sh run
user=tomcat
